# VAULT/_GARDEN/FOREST/_SCHEDULE_FREEZE_2026-01-15_13-48-54/genesis-forest-tick.timer
[Unit]
Description=Run GENESIS Forest Tick weekly

[Timer]
OnCalendar=Mon *-*-* 12:31:00
Persistent=true

[Install]
WantedBy=timers.target
