# _TRACES/POE/_CLEAN_CORE_2026-01-15_12-53-38/systemd/genesis-poe-heartbeat.service
[Unit]
Description=GENESIS POE Heartbeat (minimal trace)

[Service]
Type=oneshot
ExecStart=/home/arc/GENESIS/bin/genesis_poe_heartbeat.sh
