# _TRACES/POE/_CLEAN_CORE_2026-01-15_12-53-38/systemd/genesis-poe-proof-week.timer
[Unit]
Description=Run GENESIS POE Weekly Proof (export visible artifact)

[Timer]
OnCalendar=Mon *-*-* 12:25:00
Persistent=true

[Install]
WantedBy=timers.target
