# _GARDEN/FOREST/_SCHEDULE_FREEZE_2026-01-15_13-48-54/genesis-trader-card-poe.service
[Unit]
Description=GENESIS Trader Card (POE) - weekly offer card

[Service]
Type=oneshot
ExecStart=/home/arc/GENESIS/bin/genesis_trader_card_poe.sh
