# _GARDEN/FOREST/_SCHEDULE_FREEZE_2026-01-15_13-48-54/genesis-hybrid-tick.service
[Unit]
Description=GENESIS Hybrid Tick (weekly hybrid seed)

[Service]
Type=oneshot
ExecStart=/home/arc/GENESIS/bin/genesis_hybrid_tick.sh
