# _GARDEN/FOREST/_SCHEDULE_FREEZE_2026-01-15_13-48-54/genesis-poe-proof-week.service
[Unit]
Description=GENESIS POE Weekly Proof (core visible artifact)

[Service]
Type=oneshot
ExecStart=/home/arc/GENESIS/bin/genesis_poe_proof_week.sh
