# _GARDEN/FOREST/_SCHEDULE_FREEZE_2026-01-15_13-48-54/genesis-hybrid-tick.timer
[Unit]
Description=Run GENESIS Hybrid Tick weekly

[Timer]
OnCalendar=Mon *-*-* 12:33:00
Persistent=true

[Install]
WantedBy=timers.target
