# _GARDEN/FOREST/_SCHEDULE_FREEZE_2026-01-15_13-48-54/genesis-forest-tick.service
[Unit]
Description=GENESIS Forest Tick (biomes + crossbreed + vitrina)

[Service]
Type=oneshot
ExecStart=/home/arc/GENESIS/bin/genesis_forest_tick.sh
