# _GARDEN/FOREST/_SCHEDULE_FREEZE_2026-01-15_13-48-54/genesis-garden-harvest.service
[Unit]
Description=GENESIS Garden Harvest (weekly fruits+seeds+vitrina)

[Service]
Type=oneshot
ExecStart=/home/arc/GENESIS/bin/genesis_garden_harvest.sh
