# _GARDEN/FOREST/_SCHEDULE_FREEZE_2026-01-15_13-48-54/genesis-trader-card-poe.timer
[Unit]
Description=Run GENESIS Trader Card (POE) weekly

[Timer]
OnCalendar=Mon *-*-* 12:27:00
Persistent=true

[Install]
WantedBy=timers.target
